# scripts/com.drumgen.dailygitpush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.drumgen.dailygitpush</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/qa_m2/Documents/Cursor AI/DrumGen Scorer/scripts/daily_git_push.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/qa_m2/Documents/Cursor AI/DrumGen Scorer/scripts/git_push.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/qa_m2/Documents/Cursor AI/DrumGen Scorer/scripts/git_push_error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
